# notify-listener.service
[Unit]
Description=IRSSI notify listener

[Service]
ExecStart=/usr/bin/python2 /usr/bin/notossh -f start

[Install]
WantedBy=default.target
