# sleep.target
[Unit]
Description=Sleep
DefaultDependencies=no
StopWhenUnneeded=yes
