# xinit.target
[Unit]
Description=Xinitrc Stuff
Wants=xmonad.target

[Install]
Alias=default.target
